# etc/cron.d/cleanup-cgroup
31 3 * * * root /usr/local/bin/dgx-cgroup-old-cleanup >> /tmp/cron.dgx-cgroup-old-cleanup 2>&1
